# docker-rocky8.6-pagure-5.13.3/src/files/systemd/pagure_mirror_project_in.timer
[Unit]
Description=Update in mirrored project

[Timer]
OnCalendar=minutely
Persistent=true
AccuracySec=1min
RandomizedDelaySec=15

[Install]
WantedBy=timers.target
